# deploy/pv-materiel.service
[Unit]
Description=PV Matériel Loué - Centrale Lyon Conseil
After=network.target

[Service]
Type=notify
User=www-data
Group=www-data
WorkingDirectory=/var/www/pv-materiel
Environment="PATH=/var/www/pv-materiel/venv/bin"
Environment="FLASK_ENV=production"
ExecStart=/var/www/pv-materiel/venv/bin/gunicorn -c gunicorn_config.py server:app

# Redémarrage automatique en cas d'erreur
Restart=always
RestartSec=10

# Sécurité
PrivateTmp=true
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/pv-materiel/saved_pv /var/www/pv-materiel/config /var/log/pv-materiel

# Limites de ressources
LimitNOFILE=4096
LimitNPROC=512

[Install]
WantedBy=multi-user.target
